# nvsync.service
[Unit]
Description=nvsync VRR Profile Auto-Switching Daemon
Documentation=https://github.com/ghostkellz/nvsync
After=graphical-session.target

[Service]
Type=simple
ExecStart=/usr/bin/nvsync daemon
ExecReload=/bin/kill -HUP $MAINPID
Restart=on-failure
RestartSec=5

# Environment
Environment=XDG_CONFIG_HOME=%h/.config

[Install]
WantedBy=graphical-session.target
